# scripts/render-worker/percho-render-worker.service
[Unit]
Description=Percho render worker (agent-generated home tour videos)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/Percho
ExecStart=/usr/bin/python3 scripts/render-worker/worker.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/percho-render-worker.log
StandardError=append:/var/log/percho-render-worker.log

[Install]
WantedBy=multi-user.target
